# deploy/app-host/systemd/atlas-apphost@.service
[Unit]
Description=Atlas AppHost Instance (%i)
After=network.target
Wants=network-online.target

[Service]
Type=notify
User=atlas
Group=atlas
WorkingDirectory=/opt/atlas/instances/%i/current/backend
ExecStart=/opt/atlas/instances/%i/current/backend/Atlas.AppHost
Restart=on-failure
RestartSec=5s
TimeoutStartSec=30
TimeoutStopSec=15

# Environment
Environment=ASPNETCORE_ENVIRONMENT=Production
Environment=ASPNETCORE_URLS=http://+:${PORT}
Environment=AppInstance__ConfigPath=/opt/atlas/instances/%i/config/appinstance.json
Environment=APP_INSTANCE_HOME=/opt/atlas/instances/%i
EnvironmentFile=-/opt/atlas/instances/%i/config/env

# Resource limits
LimitNOFILE=65535
MemoryMax=512M

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/atlas/instances/%i
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=atlas-apphost-%i

[Install]
WantedBy=multi-user.target
