# deploy/platform-host/systemd/atlas-platformhost.service
[Unit]
Description=Atlas PlatformHost
After=network.target
Wants=network-online.target

[Service]
Type=notify
User=atlas
Group=atlas
WorkingDirectory=/opt/atlas/platform-host
ExecStart=/opt/atlas/platform-host/Atlas.PlatformHost
Restart=on-failure
RestartSec=5s
TimeoutStartSec=60
TimeoutStopSec=30

# Environment
Environment=ASPNETCORE_ENVIRONMENT=Production
Environment=ASPNETCORE_URLS=http://+:5001
EnvironmentFile=-/opt/atlas/platform-host/config/env

# Resource limits
LimitNOFILE=65535
MemoryMax=1G

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/atlas
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=atlas-platformhost

[Install]
WantedBy=multi-user.target
