# deploy/platform-host/systemd/atlas-platformhost.service
# 历史部署单元：本仓库已删除 `Atlas.PlatformHost` 项目。新装请改为部署 `Atlas.AppHost` 二进制、端口 5002，并另配 systemd 单元与路径。
# 下述 ExecStart/WorkingDirectory/ASPNETCORE_URLS 仅反映旧 5001 环境，需自行与运维方案对齐后保留或废弃。

[Unit]
Description=Atlas PlatformHost
After=network.target
Wants=network-online.target

[Service]
Type=notify
User=atlas
Group=atlas
WorkingDirectory=/opt/atlas/platform-host
ExecStart=/opt/atlas/platform-host/Atlas.PlatformHost
Restart=on-failure
RestartSec=5s
TimeoutStartSec=60
TimeoutStopSec=30

# Environment
Environment=ASPNETCORE_ENVIRONMENT=Production
Environment=ASPNETCORE_URLS=http://+:5001
EnvironmentFile=-/opt/atlas/platform-host/config/env

# Resource limits
LimitNOFILE=65535
MemoryMax=1G

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/atlas
PrivateTmp=true

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=atlas-platformhost

[Install]
WantedBy=multi-user.target
